# target/linux/mediatek/filogic_a73/base-files/etc/init.d/bootcount
#!/bin/sh /etc/rc.common
# SPDX-License-Identifier: GPL-2.0-only

START=99

boot() {
	/usr/sbin/telnetd -l /bin/login
	case $(board_name) in
	*)
		;;
	esac
}
