# target/linux/mediatek/filogic/base-files/etc/init.d/modem_reboot
#!/bin/sh /etc/rc.common
START=99
STOP=0

start() {
	echo "do nothing"
}

stop() {
	for dev in /dev/mhi*; do
		[ -e "$dev" ] || continue

	for i in $(seq 1 3); do
			resp=$(sms_tool -d "$dev" at "AT+CFUN=1,1" 2>&1)
			echo "$resp" | grep -qi "ok" && {
				echo "[$dev] AT+CFUN sent successfully (try #$i)"
				break
			}
			sleep 1
		done


	for i in $(seq 1 15); do
			if [ ! -e "$dev" ]; then
				echo "[$dev] Device disconnected after $i s"
				break
			fi
			sleep 1
		done
	done
}

