# target/linux/mediatek/filogic_a73/base-files/etc/uci-defaults/05_fix-compat-version
. /lib/functions.sh

case "$(board_name)" in
	*)
		;;
esac

exit 0
