# package/mtk/applications/5g-modem/luci-app-c5800-688-wtmodem/meta/wtmodem.list
/usr/bin/pdul.lua
/usr/lib/lua/luci/model/cbi/modem.lua
/usr/bin/led.sh
/usr/share/modem/RM500X
/usr/lib/libbdinfo.so.3.1
/usr/lib/libbdinfo.so.3.3
/usr/share/modem/ipcheck.sh
/usr/lib/lua/luci/view/zmode/net_status_RM520.htm
/usr/bin/full_ipv6.run
/usr/share/modem/moimei
/usr/share/modem/500U.sh
/usr/share/modem/enableipv6.sh
/etc/config/modem
/usr/share/modem/smsc.sh
/usr/bin/5700_get_zbjh.sh
/usr/share/modem/rsmsc.sh
/usr/lib/lua/luci/view/zmode/at.htm
/usr/lib/lua/luci/view/zmode/smsc.htm
/usr/bin/bdinfo1
/usr/share/modem/Quectel
/usr/bin/half_ipv6.run
/usr/bin/setsmstitle.sh
/usr/share/modem/autofreqlock.sh
/usr/lib/lua/luci/view/zmode/net_status_MT5700M.htm
/usr/share/modem/mt5700m.sh
/usr/bin/turn_off_ipv6.run
/usr/bin/BLUE4WAN.sh
/usr/share/modem/rm520n.sh
/usr/share/modem/delatcmd.sh
/usr/lib/lua/luci/view/zmode/at2.htm
/usr/share/modem/atcmd.sh
/usr/bin/smstrun.py
/usr/bin/bdinfoemmc
/usr/lib/libbdinfo.so
/usr/share/modem/mopdu
/usr/bin/l2tp-sim.sh
/usr/bin/QFirehose
/usr/lib/lua/luci/controller/modem.lua
/etc/config/apninfo
/usr/lib/lua/luci/view/zmode/net_status.htm
/usr/bin/smstrun.sh
/etc/init.d/modeminit
/usr/bin/cgsys1.sh
/usr/lib/lua/luci/view/zmode/net_status_RM500U.htm
/usr/share/modem/pingCheck.sh
/usr/bin/a.sh
/usr/share/modem/zinfo.sh
/usr/bin/l2tpsim.sh
/usr/bin/disl2tp-sim.sh
/usr/bin/dxzf.sh
/usr/bin/httpapi.py
/usr/bin/sll.sh
/usr/bin/rsrp2rssi
/usr/bin/smstrun-title.conf
/usr/share/modem/netmodeled.sh
/usr/bin/setppstoken.sh
/usr/bin/pdu_decoder
/usr/share/modem/zinfo_mt5700.sh
/usr/bin/sendat
